# scripts/conf/supervisord/supervisord.d/httpd.ini
[program:apache]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=true
#user=root
#numprocs=10
redirect_stderr=true
stdout_logfile=/var/log/httpd.log
