# scripts/conf/supervisord/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
#user=root
#numprocs=10
redirect_stderr=true
stdout_logfile=/var/log/php-fpm.log
